# xinetd.d/daytime-udp
service daytime-udp
{
	disable		= yes
	type		= INTERNAL
	id		= daytime-dgram
	socket_type	= dgram
	wait		= yes
	user		= root
	groups		= yes
	flags		= REUSE
}
