# xinetd.d/daytime
service daytime
{
	disable		= yes
	type		= INTERNAL
	id		= daytime-stream
	socket_type	= stream
	wait		= no
	user		= root
	groups		= yes
	flags		= REUSE
}
