# xinetd.d/echo
service echo
{
	disable		= yes
	type		= INTERNAL
	id		= echo-stream
	socket_type	= stream
	wait		= no
	user		= root
	groups		= yes
	flags		= REUSE
}
