# xinetd.d/finger
service finger
{
	disable		= yes
	socket_type	= stream
	wait		= no
	user		= nobody
	server		= /usr/libexec/fingerd
	server_args	= -s
	groups		= yes
	flags		= REUSE
}
