# xinetd.d/ftp
service ftp
{
	disable		= yes
	socket_type	= stream
	wait		= no
	user		= root
	server		= /usr/libexec/ftpd
	server_args	= -l
	groups		= yes
	flags		= REUSE IPv6
}
