# xinetd.d/ntalk
service ntalk
{
	disable		= yes
	socket_type	= dgram
	wait		= yes
	user		= root
	server		= /usr/libexec/ntalkd
	groups		= yes
	flags		= REUSE
}
