# xinetd.d/telnet
service telnet
{
	disable		= yes
	socket_type	= stream
	wait		= no
	user		= root
	server		= /usr/libexec/telnetd
	groups		= yes
	flags		= REUSE
}
