# xinetd.d/tftp
service tftp
{
	cps		= 200 5
	disable		= yes
	socket_type	= dgram
	wait		= yes
	user		= nobody
	server		= /usr/libexec/tftpd
	server_args	= /private/tftpboot
	groups		= yes
	flags		= REUSE
}
